# supervisord-serf-service-registry.conf
[program:serf-service-registry]
command=/start-serf-service-registry.sh
numprocs=1
autostart=true
autorestart=true

[program:serf-join]
command=/serf-join.sh
autorestart=false
